# rootfs/etc/s6-overlay/s6-rc.d/scheduler/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# WieZorgt scheduler service
# ==============================================================================
bashio::log.info "Starting scheduler service"

cd /usr/bin/ai
# Run daemon
exec python3 scheduler.py
